# deploy/systemd/arkeep-agent.service
[Unit]
Description=Arkeep Backup Agent
Documentation=https://arkeep.io/docs
After=network-online.target
Wants=network-online.target
# If Docker is present, wait for it so volume discovery works on startup.
After=docker.service

[Service]
Type=simple
User=arkeep-agent
Group=arkeep-agent

ExecStart=/usr/local/bin/arkeep-agent \
    --server-addr=${ARKEEP_SERVER_ADDR} \
    --state-dir=/var/lib/arkeep-agent \
    --log-level=${ARKEEP_LOG_LEVEL:-info}

# Environment file — copy to /etc/arkeep/agent.env and fill in the values.
# Minimum required:
#   ARKEEP_SERVER_ADDR=arkeep.example.com:9090
#   ARKEEP_AGENT_SECRET=<shared secret — must match the server>
#
# Auto-enrollment (default): on first startup the agent calls the server's HTTP
# API (port 8080, derived automatically) to obtain a signed client certificate.
# The cert is stored in /var/lib/arkeep-agent and reused on every restart.
#
# Override the enrollment HTTP address when the server is behind a reverse
# proxy (Traefik, Nginx, Caddy) — the derived http://<host>:8080 won't work:
#   ARKEEP_SERVER_HTTP_ADDR=https://arkeep.example.com
# If port 8080 is exposed directly, the default is sufficient:
#   ARKEEP_SERVER_HTTP_ADDR=http://arkeep.example.com:8080
#
# Only needed when the server uses an external cert from a non-system CA:
#   ARKEEP_GRPC_TLS_CA=/etc/arkeep/ca.crt
#
# Disable TLS (same-machine only — must also set ARKEEP_GRPC_INSECURE on server):
#   ARKEEP_GRPC_INSECURE=true
EnvironmentFile=/etc/arkeep/agent.env

Restart=on-failure
RestartSec=10s
TimeoutStopSec=30s

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
# Add the agent user to the docker group for Docker socket access.
# Required for automatic container and volume discovery.
# If Docker is not installed you can remove this line.
SupplementaryGroups=docker

StateDirectory=arkeep-agent

[Install]
WantedBy=multi-user.target
